# supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=info
logfile=/dev/stdout 
logfile_maxbytes=0  
pidfile=/tmp/supervisord.pid

[program:app]
command=/app/app

[program:ssr]
command=bun run ssr/ssr.js